# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root

[program:caddy]
command=/app/scripts/caddy-start.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
environment=XDG_DATA_HOME="/data/caddy",GC_CADDY_CONFIG_PATH="%(ENV_GC_CADDY_CONFIG_PATH)s"

[program:wireguard]
command=/app/scripts/wg-wrapper.sh
autostart=true
autorestart=unexpected
exitcodes=0
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=3
priority=20
environment=GC_WG_INTERFACE="%(ENV_GC_WG_INTERFACE)s"

; dnsmasq lifecycle is owned by scripts/wg-wrapper.sh — it can only
; bind 10.8.0.1:53 once wg-quick has brought wg0 up, and racing it
; here against [program:wireguard] left supervisord in a permanent
; FATAL retry loop after wg-wrapper killed the early instance.
;
; Split-horizon DNS still hijacks GC_API_HOST → GC_WG_GATEWAY_IP for
; VPN peers; the dnsmasq-process just starts later in the boot.

[program:node]
command=node /app/src/server.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
startretries=5
startsecs=3
