# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root

[program:caddy]
command=caddy run --config /app/config/Caddyfile --adapter caddyfile
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
environment=XDG_DATA_HOME="/data/caddy"

[program:wireguard]
command=/bin/sh -c "wg-quick up %(ENV_GC_WG_INTERFACE)s 2>&1 && sleep infinity"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=0
priority=20
environment=GC_WG_INTERFACE="%(ENV_GC_WG_INTERFACE)s"

[program:node]
command=node /app/src/server.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
startretries=5
startsecs=3
